# archiso/kde/neco-pkg/pkg/necoarch/etc/profile.d/guranyu.sh
#!/bin/sh
cvlc --play-and-exit guranyu.mp3 &>/dev/null &
